# lib/systemd/system/raspiaprs.service
[Unit]
Description=APRS Metrics
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=45
ProtectSystem=true
ProtectHome=true
PermissionsStartOnly=true
ExecStart=/usr/local/bin/raspiaprs

[Install]
Alias=raspiaprs.service
WantedBy=multi-user.target
