# raspiaprs.service
[Unit]
Description=RasPiAPRS Monitoring Service
After=network.target gpsd.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/RasPiAPRS
ExecStart=/bin/bash /opt/RasPiAPRS/main.sh
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT
KillMode=process
Restart=always
RestartSec=20
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target